# sales-financial-backend.service
[Unit]
Description=Sales Financial Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/your/project
EnvironmentFile=/path/to/your/project/.env.production
ExecStart=/usr/bin/node dist/__create/index.js
Restart=always

[Install]
WantedBy=multi-user.target
